# roles/eosc_vd_portal/templates/tomcat.service.j2
[Unit]
Description=Tomcat

[Service]
Type=forking
ExecStart={{tomcat_dir}}/bin/startup.sh
KillMode=process
User=root
Group=root

[Install]
WantedBy=multi-user.target